# docker/postgres/crontab
* * * * * /var/www/data/database/dump.sh